# A.service
[Unit]
Description=Service A

[Service]
Type=notify
ExecStart=/tmp/a.out

[Install]
WantedBy=multi-user.target

